# nextcloud/cron/nextcloud-cron.timer
[Unit]
Description=Nextcloud cron timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
Unit=nextcloud-cron.service

[Install]
WantedBy=timers.target
